# supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
autorestart=true
startsecs=10
startretries=3
command=/usr/local/bin/uwsgi --ini /web_app/uwsgi_config.ini
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
